# etc/crontab
30 * * * * /usr/local/sbin/snapshot-hourly
15 3 * * * /usr/local/sbin/snapshot-daily
0 5 * * 0 /usr/local/sbin/update-server
